# docs/etc/init/app-web-1.conf
start on starting app-web
stop on stopping app-web
respawn

env PORT=5000

setuid app

chdir /Users/apple/Projects/stream_alert

exec bundle exec thin start -p $PORT
