# docs/etc/init/app-memory-1.conf
start on starting app-memory
stop on stopping app-memory
respawn

env PORT=5300

setuid app

chdir /Users/apple/Projects/stream_alert

exec redis-server
