# docs/etc/init/app-worker-1.conf
start on starting app-worker
stop on stopping app-worker
respawn

env PORT=5100

setuid app

chdir /Users/apple/Projects/stream_alert

exec bundle exec sidekiq
