# docs/etc/init/app-search-1.conf
start on starting app-search
stop on stopping app-search
respawn

env PORT=5200

setuid app

chdir /Users/apple/Projects/stream_alert

exec elasticsearch
